# scripts/com.example.job-search-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd template for weekday daily job search.
  Before use:
    1. Replace __REPO_ROOT__ with the absolute path to your job-search clone
    2. Replace com.example.job-search-daily with your reverse-DNS label if desired
    3. cp scripts/com.example.job-search-daily.plist ~/Library/LaunchAgents/
    4. launchctl load ~/Library/LaunchAgents/com.example.job-search-daily.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.job-search-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/run-daily-search.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/data/logs/launchd-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/data/logs/launchd-stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
